# appsflyer-dashboard.service
[Unit]
Description=AppsFlyer Dashboard Service
After=network.target redis.service

[Service]
Type=simple
User=or
Group=or
WorkingDirectory=/home/or/appsflyer-dash-2025-live
Environment="PATH=/home/or/appsflyer-dash-2025-live/venv/bin"
Environment="FLASK_ENV=production"
Environment="FLASK_DEBUG=0"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/or/appsflyer-dash-2025-live/venv/bin/gunicorn --workers 4 --bind 0.0.0.0:5000 --timeout 3600 --log-level info --access-logfile /home/or/appsflyer-dash-2025-live/logs/access.log --error-logfile /home/or/appsflyer-dash-2025-live/logs/error.log app:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 